# logrotate/logrotate.cron
MAILTO=its-iam@lists.hawaii.edu
@daily /home/grouper/bin/logrotate-wrapper.sh